# deploy/tianyi/huge-catering.service
[Unit]
Description=Huge Catering WeChat Article Generator
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/projects/HUGEGONGZHONGHAO
EnvironmentFile=-/opt/projects/HUGEGONGZHONGHAO/.env
Environment=PYTHONPATH=/opt/projects/HUGEGONGZHONGHAO/src
ExecStart=/opt/projects/HUGEGONGZHONGHAO/.venv/bin/gunicorn --workers 2 --bind 0.0.0.0:8766 huge_catering.wsgi:app
Restart=always
RestartSec=5
TimeoutStartSec=30
KillSignal=SIGQUIT

[Install]
WantedBy=multi-user.target
